# tarantula/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/usr/bin/mysqld_safe

[program:httpd]
command=/usr/sbin/apachectl -D FOREGROUND

