# omero-deploy/supervisord/omero.conf
[program:omero-setup]
command=/omero-setup
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=900
