# omero-deploy/supervisord/ssh.conf
[program:ssh]
command=/usr/sbin/sshd -D
autorestart=true
priority=100
