# omero-server-prereq/supervisord/postgresql.conf
[program:postgresql]
command=/postgresql
autostart=true
autorestart=false
stopsignal=QUIT
stopwaitsecs=120
priority=200
